# net-misc/janus/files/janus-initd
#!/sbin/runscript
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

JANUS_BIN="${JANUS_BIN:-/usr/bin/janus}"

depend() {
  need net
}


janus_svcname() {
        local ebextra=
        case "${SVCNAME}" in
                janus*) ;;
                *) ebextra=" ``(janus)" ;;
        esac
        echo "${SVCNAME}${ebextra}"
}

do_setlimits() {
        einfo "setting general limits"
        ulimit -c unlimited
        ulimit -n 64000
}

# /usr/bin/janus -F${JANUS_CNF_PATH} ${JANUS_OPTS} -N -L ${janus_log_file} -p ${pid_file}
start() {
    ebegin "Starting $(janus_svcname)"

    JANUS_PID_FILE="/var/run/${SVCNAME}.pid"
    JANUS_CNF_PATH="${JANUS_CNF_PATH:-/etc/${SVCNAME}}"
    JANUS_LOG_FILE="/var/log/janus/${SVCNAME}.log"

    if [ ! -r "${JANUS_CNF_PATH}/janus.cfg" ] ; then
            eerror "Cannot read the configuration file \`${JANUS_CNF_PATH}/janus.cfg'"
            return 1
    fi

    # set the limits
    do_setlimits


    start-stop-daemon --wait ${JANUS_RESTART_WAIT} \
    --background --start --exec \
    ${JANUS_BIN} \
    --make-pidfile --pidfile ${JANUS_PID_FILE} \
    -- -F ${JANUS_CNF_PATH} ${JANUS_OPTS} -N -L ${JANUS_LOG_FILE}
    eend $?
}


stop() {
    ebegin "Stopping $(janus_svcname)"

    JANUS_PID_FILE="/var/run/${SVCNAME}.pid"
    JANUS_CNF_PATH="${JANUS_CNF_PATH:-/etc/${SVCNAME}}"
    JANUS_LOG_FILE="/var/log/janus/${SVCNAME}.log"

    start-stop-daemon --wait ${JANUS_RESTART_WAIT} \
    --stop --exec ${JANUS_BIN} \
    --pidfile ${JANUS_PID_FILE}
    eend $?
}

reload() {
    ebegin "Reloading myApp"
    start-stop-daemon --exec ${JANUS_BIN} \
    --pidfile ${JANUS_PID_FILE} \
    -s 1
    eend $?
}
